# systemd/piun.service
[Unit]
Description=Podman Image Update Notifier

[Service]
Type=oneshot
ExecStart=%h/.local/bin/piun
