# systemd/piun.timer
[Unit]
Description=Podman Image Update Notifier

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
