# systemd/piun.service
[Unit]
Description=Podman Image Update Notifier

[Service]
Type=simple
ExecStart=%h/.local/bin/piun

[Install]
WantedBy=default.target
