# systemd/piun.timer
[Unit]
Description=Podman Image Update Notifier

[Timer]
OnUnitInactiveSec=12h
RemainAfterElapse=no

[Install]
WantedBy=timers.target
